# systemd/pm-second-brain-daily-brief.timer.tmpl
[Unit]
Description={{UNIT}}

[Timer]
OnCalendar={{ON_CALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
