# systemd/pm-second-brain-daily-brief.service.tmpl
[Unit]
Description={{UNIT}}

[Service]
Type=oneshot
ExecStart={{EXEC_START}}
