# launchd/com.pm-second-brain.daily-brief.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
<dict>
    <key>Label</key><string>{{LABEL}}</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{PROGRAM}}</string>
        <string>run</string>
        <string>--skill</string>
        <string>pm-workflow.pm-daily-brief</string>
        <string>--headless</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <!-- generated by pm_second_brain.cron.render_launchd_plist from {{CRON_EXPR}} -->
        {{INTERVALS}}
    </array>
    <key>StandardOutPath</key><string>{{LOG_PATH}}</string>
    <key>StandardErrorPath</key><string>{{LOG_PATH}}</string>
</dict>
</plist>
